# cde/contrib/systemd/getty@.service.d/00-cde-fix-utmp.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \\u' --noclear %I $TERM
