# cde/contrib/systemd/getty@tty1.service.d/00-cde-no-reset.conf
[Service]
TTYReset=no
TTYVTDisallocate=no
